# srvouve/etc/systemd/system/srvouve.service
[Unit]
Description=srvouve

[Service]
ExecStart=/usr/bin/srvouve

[Install]
WantedBy=multi-user.target
